# deploy/systemd/yagocrawler.service
[Unit]
Description=yago crawler worker
Documentation=https://github.com/D4rk4/yago
After=network-online.target yago-node.service
Wants=network-online.target

[Service]
Type=simple
User=yago
Group=yago
ExecStart=/opt/yago/bin/yagocrawler
EnvironmentFile=/opt/yago/etc/yagocrawler.env
WorkingDirectory=/opt/yago
Environment=HOME=/opt/yago/data/crawler
Restart=on-failure
RestartSec=5s

# The service, not a container, is the isolation boundary on bare metal. Keeping
# NoNewPrivileges requires the crawler's default YAGOCRAWLER_BROWSER_SANDBOX=false
# (Firefox cannot start its content sandbox under NoNewPrivileges); see deploy/README.md.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/opt/yago/data
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
LockPersonality=yes

MemoryHigh=60%
MemoryMax=85%
TasksMax=4096
CPUWeight=50

Nice=5
IOSchedulingClass=best-effort
IOSchedulingPriority=6
IOWeight=50

[Install]
WantedBy=multi-user.target
